# dot_config/systemd/user/restic-home.timer.tmpl
[Unit]
Description=Backs up Home Folder
Requires=restic-home.service

[Timer]
Unit=restic-home.service
OnCalendar=*:00:00/5

[Install]
WantedBy=timers.target
