# private_Library/LaunchAgents/local.bluebubbles.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.bluebubbles.tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/andrew.burns/bin/blue-bubbles-loop.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/andrew.burns</string>
    <key>PATH</key>
    <string>/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/andrew.burns</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/andrew.burns/Library/Logs/blue-bubbles-loop.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/andrew.burns/Library/Logs/blue-bubbles-loop.log</string>
</dict>
</plist>
