# direct/rotate-wallpaper.cron
# Rotate wezterm wallpaper
33 06-23 * * * ~/bin/rotate-wezterm-wallpaper
